# supervisord.conf
[program:node]
command = node src/nasjonalturbase.api.js
directory = /home/dotcloud/current
