# supervisord.conf
[program:node]
command = coffee coffee/nasjonalturbase.api.coffee
directory = /home/dotcloud/current
